# docs/systemd/afnsec-ssh-gate-purge.service.example
[Unit]
Description=AFNSec SSH Gate cache purge
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/afnsec-ssh-gate --purge
User=root
Group=root
